# logtool.service
[Unit]
Description=logtool 系统日志分析守护进程
Documentation=https://github.com/your-repo/logtool
After=systemd-journald.service
Wants=systemd-journald.service

[Service]
Type=simple
# 与 deb 安装路径保持一致（/usr/bin/logtool-daemon）。
ExecStart=/usr/bin/logtool-daemon
ExecStartPre=/bin/rm -f /run/logtool.sock
Restart=on-failure
RestartSec=5

# 安全加固
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/
ReadWritePaths=/run
NoNewPrivileges=true
PrivateTmp=true

# Socket 由守护进程创建为 0660，推荐 root:logtool。
# 使用 logtool CLI 的普通用户需加入 logtool 组。

[Install]
WantedBy=multi-user.target
